# syswatch-logger.service
[Unit]
Description=syswatch background metrics logger
After=network.target

[Service]
Type=simple
# User= is rewritten to the installing user by install-syswatch.sh; the value
# below is just the default fallback used when SUDO_USER is unset.
User=lukas
ExecStart=/usr/bin/python3 /usr/local/lib/syswatch/syswatch-logger.py
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
